# bandwidth-monitor.service
[Unit]
Description=Bandwidth Monitor
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/bandwidth-monitor
EnvironmentFile=/opt/bandwidth-monitor/.env
ExecStart=/opt/bandwidth-monitor/bandwidth-monitor
Restart=on-failure
RestartSec=5
AmbientCapabilities=CAP_NET_RAW
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadOnlyPaths=/opt/bandwidth-monitor
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
